# config/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command = /usr/bin/uwsgi --ini /apps/storefront/config/storefront.ini

[program:nginx]
command = /usr/sbin/nginx -g 'daemon off; user nginx nginx;'

[program:cron]
command = cron -f -L 15
autostart=true
autorestart=true
